# metropolitan-bot.service
[Unit]
Description=Metropolitan Expo Attendance Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/metrpolitanbot
Environment=PATH=/path/to/your/metrpolitanbot/venv/bin
ExecStart=/path/to/your/metrpolitanbot/venv/bin/python attendance_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=metropolitan-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/your/metrpolitanbot/logs

[Install]
WantedBy=multi-user.target
